# supervisor-config.ini
[program:timao-main]
command=/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin/python -m uvicorn server.main:app --host 0.0.0.0 --port 8000
directory=/www/wwwroot/wwwroot/timao-douyin-live-manager
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/www/wwwroot/wwwroot/timao-douyin-live-manager/logs/main.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin:%(ENV_PATH)s"

[program:timao-streamcap]
command=/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin/python -m uvicorn server.streamcap_service:app --host 0.0.0.0 --port 8001
directory=/www/wwwroot/wwwroot/timao-douyin-live-manager
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/www/wwwroot/wwwroot/timao-douyin-live-manager/logs/streamcap.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin:%(ENV_PATH)s"

[program:timao-douyin]
command=/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin/python -m uvicorn server.douyin_service:app --host 0.0.0.0 --port 8002
directory=/www/wwwroot/wwwroot/timao-douyin-live-manager
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/www/wwwroot/wwwroot/timao-douyin-live-manager/logs/douyin.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/www/wwwroot/wwwroot/timao-douyin-live-manager/venv/bin:%(ENV_PATH)s"

[group:timao-services]
programs=timao-main,timao-streamcap,timao-douyin
priority=999

